# build/cmake/debian/postinst.in
#!/bin/sh -e

echo "Running $0 $@"

case "$1" in
    configure)
        if type systemctl >/dev/null 2>&1
        then
            echo "Adding system service: @SERVICE_UNITS@"
            systemctl enable @SERVICE_UNIT_PATHS@
            systemctl start @SERVICE_UNITS@
        else
            echo "SystemD is not available; skipping service integration"
        fi
        ;;
esac
