# build/debian/prerm.in
#!/bin/sh -e

echo "Running $0 $@"

case "$1" in
    upgrade)
    ;;

    remove)
        if type systemctl >/dev/null 2>&1
        then
            echo "Removing system service: @SERVICE_UNITS@"
            systemctl disable @SERVICE_UNITS@ || true
            systemctl stop @SERVICE_UNITS@ || true
        else
            echo "SystemD is not available; skipping service integration"
        fi
        ;;
esac
